# classes/ntut_embedded_systems/review/zuvio/06_email.crontab
# Schedule an email to kiki on 14/10/2026 only for that year
# Crontab entry

0 0 14 10 * [ $(date +\%Y) -eq 2026 ] && echo "I am sending an email to you." | mail -s "Scheduled Message for 2026" kiki